# .config/systemd/user/barriers.service
[Unit]
Description=Barrier Server

[Service]
ExecStart=/usr/bin/barriers -f --disable-client-cert-checking

[Install]
WantedBy=default.target
