# .config/systemd/user/edges.service
[Unit]
Description=Edges

[Service]
ExecStart=/usr/bin/edges --bottom eww_mouse.sh

[Install]
WantedBy=default.target
