# .config/systemd/user/eww.service
[Unit]
Description=Eww Daemon

[Service]
ExecStart=/usr/bin/eww --no-daemonize daemon

[Install]
WantedBy=default.target
