# .config/systemd/user/flashfocus.service
[Unit]
Description=Flashfocus

[Service]
ExecStart=/usr/bin/flashfocus
Restart=on-failure
RestartSec=3s

[Install]
WantedBy=default.target
