# .config/systemd/user/flashfocus.service
[Unit]
Description=Flashfocus

[Service]
ExecStart=/usr/bin/flashfocus

[Install]
WantedBy=default.target
