# .config/systemd/user/greenclip.service
[Unit]
Description=Greenclip Daemon

[Service]
ExecStart=/usr/bin/greenclip daemon

[Install]
WantedBy=default.target
