# .config/systemd/user/light-locker.service
[Unit]
Description=Light Locker Daemon

[Service]
ExecStart=/usr/bin/light-locker

[Install]
WantedBy=default.target
