# .config/systemd/user/picom.service
[Unit]
Description=Picom

[Service]
ExecStart=/usr/bin/picom --experimental-backends

[Install]
WantedBy=default.target
