# .config/systemd/user/xbanish.service
[Unit]
Description=xbanish

[Service]
ExecStart=/usr/bin/xbanish

[Install]
WantedBy=default.target
